# buildroot/package/snapdog-rauc/rauc-mark-good.service
[Unit]
Description=Mark RAUC slot as good after successful boot
# Wants (not Requires): always run so rauc-commit's health gate decides whether
# to commit; a failed ctrl must not tear this unit down before the gate runs.
After=snapdog-ctrl.service
Wants=snapdog-ctrl.service
ConditionPathExists=/usr/bin/rauc

[Service]
Type=oneshot
# Marks the booted slot good AND commits/reconciles a pending tryboot trial.
ExecStart=/opt/snapdog/bin/rauc-commit

[Install]
WantedBy=multi-user.target
